# Ansible/provision/roles/bluepill/templates/etc_init_bluepill_config.conf.j2
description "Application_new Bluepill"
start on runlevel [2]
stop on runlevel [016]
expect daemon
exec bluepill load {{ config_directory }}/bluepill.rb
respawn
